# controlroom/jobs.cron
PATH=/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
0 13 * * * /Users/sam/autotv/controlroom/run.sh generate-shopping
41 12 * * * /Users/sam/autotv/controlroom/run.sh generate-homes
30 1 * * * rm /Users/sam/autotv/cookingshow2/videos/*.mp4
#0 2 * * * /Users/sam/autotv/controlroom/run.sh generate-cooking
